# archies/.config/systemd/user/bspwm-reload.service
[Unit]
Description=Reload BSPWM

[Service]
Type=oneshot
ExecStart=/bin/bash -c "/home/hubbe/.scripts/monitor_manager"
StandardOutput=journal
